# aistock-cronjob/crontab
# AI Stock 定时任务配置
# 格式: 分钟 小时 日 月 星期 用户 命令

# 获取市场指数 - 每 10 秒执行
* * * * * root cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1
* * * * * root sleep 10; cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1
* * * * * root sleep 20; cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1
* * * * * root sleep 30; cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1
* * * * * root sleep 40; cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1
* * * * * root sleep 50; cd /app && /usr/local/bin/python get_index.py >> /var/log/cron/index.log 2>&1

# 获取新闻 - 每 10 分执行
*/10 * * * * root cd /app && /usr/local/bin/python get_news.py >> /var/log/cron/news.log 2>&1

# 股票详细页定时缓存 - 每 8 分执行
*/8 * * * * root cd /app && /usr/local/bin/python cache_stock_detail_page.py >> /var/log/cron/cache_stock_detail_page.log 2>&1

# 获取实时行情数据 - 每 10 分执行
*/10 * * * * root cd /app && /usr/local/bin/python get_realtime_quotes.py >> /var/log/cron/realtime_quotes.log 2>&1

# 获取热门股票 - 每 30 分执行
*/30 * * * * root cd /app && /usr/local/bin/python get_hot_stock.py >> /var/log/cron/hot_stock.log 2>&1

# 获取标签对应股票 - 每 30 分执行
*/30 * * * * root cd /app && /usr/local/bin/python get_tags_stocks.py >> /var/log/cron/tags_stocks.log 2>&1

# 更新个股新闻 - 每 1 小时 30 分执行
30 */1 * * * root cd /app && /usr/local/bin/python get_stock_news.py >> /var/log/cron/stock_news.log 2>&1

# 获取历史行情数据 - 每日 0:00 执行
0 0 * * * root cd /app && /usr/local/bin/python get_stock_history.py >> /var/log/cron/stock_history.log 2>&1

# 更新个股 ai 评测 - 每日 0:30 执行
30 0 * * * root cd /app && /usr/local/bin/python get_stock_eva.py >> /var/log/cron/stock_eva.log 2>&1

# 获取股票列表 - 每日 2:30 执行
30 2 * * * root cd /app && /usr/local/bin/python get_stock_name.py >> /var/log/cron/stock_name.log 2>&1

# 早晨推送 - 每日 9:30 执行
30 9 * * * root cd /app && /usr/local/bin/python wechat_push_morning.py >> /var/log/cron/wechat_morning.log 2>&1

# 更新个股 ai 评测 - 每日 9:30 执行
30 9 * * * root cd /app && /usr/local/bin/python get_stock_eva.py >> /var/log/cron/stock_eva.log 2>&1

# 定时推送 - 每日 10:00、11:00、12:00、13:00、14:00、15:00、16:00 执行
0 10,11,12,13,14,15,16 * * * root cd /app && /usr/local/bin/python wechat_push.py >> /var/log/cron/wechat_push.log 2>&1

# 更新个股 ai 评测 - 每日 14:30 执行
30 14 * * * root cd /app && /usr/local/bin/python get_stock_eva.py >> /var/log/cron/stock_eva.log 2>&1

# 更新个股 ai 评测 - 每日 18:30 执行
30 18 * * * root cd /app && /usr/local/bin/python get_stock_eva.py >> /var/log/cron/stock_eva.log 2>&1

# 获取股票信息 - 每日 22:30 执行
30 22 * * * root cd /app && /usr/local/bin/python get_stock_info.py >> /var/log/cron/stock_info.log 2>&1

# 缓存清理 - 每周周一 0:00 执行
0 0 * * 1 root cd /app && /usr/local/bin/python -c "import redis; import os; r=redis.from_url(os.getenv('REDIS_BROKER_URL', 'redis://redis:6379/0')); r.flushdb()" >> /var/log/cron/cache_clean.log 2>&1

# 清理日志 - 每月 1 日 0:00 执行
0 0 1 * * root find /var/log/cron/ -type f -name "*.log" -mtime +30 -exec rm -f {} \; >> /var/log/cron/log_cleanup.log 2>&1

# 清理过期新闻 - 每月 1 日 0:00 执行
0 0 1 * * root cd /app && /usr/local/bin/python cleanup_old_news.py >> /var/log/cron/cleanup_old_news.log 2>&1

# 数据库备份 - 每 3 日 6:00 执行
0 6 */3 * * root cd /app && /usr/local/bin/python backup_mysql_minio.py >> /var/log/cron/backup_mysql.log 2>&1
